# roles/graphite/templates/supervisor.conf
[program:graphite]
command={{prefix}}/webapp/start_graphite.sh
user=root
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile = {{app_log}}/supervisor.log 
